# supervisor/nats-listener.conf
[program:nats-listener]
command=php /var/www/html/artisan nats:listen
autostart=true
autorestart=true
user=laravel
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/nats-listener.log